# systemd/weather-display.service
[Unit]
Description=WeatherDisplay - Weather + clock on Raspberry Pi
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/weather_display
WorkingDirectory=/home/pi/Documents/WeatherDisplay
Restart=on-failure
RestartSec=10
User=pi

# Pass environment from a file so secrets aren't in the unit.
EnvironmentFile=-/home/pi/Documents/WeatherDisplay/.env

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target